# a2ln.service
[Unit]
Description=Android 2 Linux Notifications Server

[Service]
ExecStart=/usr/bin/a2ln --no-pairing-server

[Install]
WantedBy=default.target
