# a2ln.service
[Unit]
Description=a2ln notification server

[Service]
ExecStart=/usr/local/bin/a2ln PORT
estart=always

[Install]
WantedBy=multi-user.target
